# upchecker/crontab
* * * * * root echo "$(date) Hello world" >> /var/log/cron.log 2>&1
${UPCHECKER_CRON} root /opt/upchecker/upchecker.sh
